# root/etc/services.d/varnish/run
#!/usr/bin/with-contenv bash

exec varnishd -a :80 -f /etc/varnish/default.vcl -F
